# launchd/com.wukong.news-analyze.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wukong.news-analyze</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/wukong/mylife/news/scripts/run_scheduled.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NEWS_RUN_MODE</key>
    <string>analyze</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/wukong/mylife/news</string>

  <!-- 18:00 local time, after A-share close. Change Hour/Minute if the desired daily analyze time changes. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <!-- The wrapper writes dated redacted logs under project-local .local/news-data/logs; launchd stdout/stderr stay quiet. -->
  <key>StandardOutPath</key>
  <string>/dev/null</string>
  <key>StandardErrorPath</key>
  <string>/dev/null</string>
</dict>
</plist>
